# examples/systemd/cisco-vk@.service
[Unit]
Description=Cisco Virtual Kubelet Provider - %i
After=network.target
Documentation=https://github.com/cisco/virtual-kubelet-cisco

[Service]
Type=simple
User=root

# Load environment from instance-specific file
EnvironmentFile=/etc/cisco-vk/%i.env

# Main executable
ExecStart=/usr/local/bin/cisco-vk

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cisco-vk-%i

[Install]
WantedBy=multi-user.target
